# overlays/users/usr/lib/systemd/system/manage-ssh-keys.service
[Unit]
Description=Manage SSH keys for the core user

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ssh-keys.sh
